# config/files/usr/etc/profile.d/theme-color.sh
COL_RESET="\033[0m"
RGB(){
	local RED=$1
	local GREEN=$2
	local BLUE=$3
	echo "\[\e[38;2;${RED};${GREEN};${BLUE}m\]"
}
COL_RED="$(RGB 255 0 0)"
COL_GREEN="$(RGB 0 255 0)"
COL_YELLOW="$(RGB 255 255 0)"
COL_BLUE="$(RGB 0 0 255)"
COL_MAGENTA="$(RGB 255 0 255)"
COL_CYAN="$(RGB 0 255 255)"
COL_WHITE="$(RGB 255 255 255)"
