# config/files/usr/etc/profile.d/theme.sh
source /etc/profile.d/.theme.$(basename $SHELL)
